# bridge/systemd/instantlink-bridge.service
[Unit]
Description=InstantLink Bridge FTP receive service
# Start the LCD/UI as soon as the system bus is available. Network, DHCP, and
# Bluetooth are wanted but intentionally not ordered before this unit; the app
# reports those subsystems as they come online.
After=dbus.service
Wants=bluetooth.service systemd-networkd.service dnsmasq.service NetworkManager.service

[Service]
Type=notify
NotifyAccess=main
User=ib
Group=ib
SupplementaryGroups=bluetooth gpio spi i2c plugdev video
WorkingDirectory=/opt/InstantLinkBridge
ExecStart=/opt/InstantLinkBridge/.venv/bin/instantlink-bridge --config /etc/InstantLinkBridge/config.toml
Restart=on-failure
RestartSec=2
WatchdogSec=30
# Stop budget: the BLE worker holds an FFI scan call that can take a few
# seconds to release after SIGTERM. The InstantLink backend wrapper has a
# 5 s in-process grace before bailing, ui.stop() then unwinds the rest of
# the task graph. 12 s was tight enough that every --restart deploy
# bumped into stop-sigterm timeout → SIGKILL → unit failed → manual
# reset-failed. 30 s leaves the FFI grace + cleanup ample room.
TimeoutStopSec=30
AmbientCapabilities=CAP_NET_BIND_SERVICE
# Wi-Fi mode switching is limited through /etc/sudoers.d/instantlink-bridge-wifi.
# Do not restrict CapabilityBoundingSet or set NoNewPrivileges=yes; sudo cannot run the
# root-owned NetworkManager helper if setuid/setgid capabilities are stripped.

[Install]
WantedBy=multi-user.target
