# bridge/systemd/instantlink-bridge-boot-splash.service
[Unit]
Description=InstantLink Bridge early LCD boot splash
# docs/plans/032: the splash is cosmetic. Previously it ordered
# `Before=basic.target instantlink-bridge.service`, which placed its ~2.4 s
# ExecStart on the bridge READY=1 critical path. With that ordering removed it
# now runs in parallel with the bridge — kept WantedBy=sysinit.target so it
# still activates as early as the rest of sysinit (visual feedback retained)
# while no longer blocking the bridge.
DefaultDependencies=no
After=local-fs.target systemd-udev-settle.service
ConditionPathExists=/opt/InstantLinkBridge/.venv/bin/instantlink-bridge-splash

[Service]
Type=oneshot
User=ib
Group=ib
SupplementaryGroups=gpio spi i2c video
WorkingDirectory=/opt/InstantLinkBridge
ExecStart=/opt/InstantLinkBridge/.venv/bin/instantlink-bridge-splash --config /etc/InstantLinkBridge/config.toml --hold 0.25
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
